# config/crontab/cron
1 8 * * *  php run.php#huang1798@skyworth.com
